# scripts/ceres-mail-ui.service
[Unit]
Description=CERES Mail UI (send onboarding package)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/ceres/Ceres
EnvironmentFile=-/etc/ceres/mail-ui.env
Environment=HOME=/root
ExecStart=/opt/ceres/Ceres/bin/ceres-mail-ui
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
